# supervisor.conf
[program:iotscape]
directory=/home/pi/Desktop/adapter/
command = screen -mDS iotscape ./iotscape.sh
autostart = true
autorestart = true
user=pi
